# crontab
*/10 * * * * /usr/bin/python3 /app/run.py
